# .config/zsh/profile/profile.d/1-ssh-agent.zsh
/home/soda/.config/zsh/profile/profile.d/.src/1-ssh-agent.zsh