# .config/zsh/profile/profile.d/99-uwsm.zsh
/home/soda/.config/zsh/profile/profile.d/.src/99-uwsm.zsh